# supervisor.conf
[program:tabby_search]
command = bash /home/vibloteam/Project/tabby_search/start_server.sh
directory = /home/vibloteam/Project/tabby_search
user = vibloteam
autostart=true
autorestart=true